# php/crontab/crontab
*/1 * * * * root date >> /var/log/php/crontab.log 2>&1
