# php/crontab/crontab
*/1 * * * * echo "ads" >> /var/log/php/crontab.log 2>&1
